# metadata-updater/aclvt-update-metadata.timer
[Unit]
Description=Update metadata for ACLVT every 30 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target